# cookbooks/cdo-apps/templates/default/init.d.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <%= @app %>
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: <%= @app %> app server
# Description:
### END INIT INFO

# LSB-compliant init.d script requirement
. /lib/lsb/init-functions

APP_ROOT=<%= @app_root %>
PID=<%= @pid_file %>
AS_USER=<%= @user %>


PASSENGER="bundle exec passenger"
PASSENGER_CONFIG="bundle exec passenger-config"
ADDRESS=127.0.0.1
PORT=<%= node['cdo-apps'][@app]['port'] %>
ENVIRONMENT=<%= @env %>
APP_NAME=<%= @app %>
SET_PATH="cd ${APP_ROOT}"
PRE_CMD=""
CMD="${SET_PATH}; ${PRE_CMD} ${PASSENGER} start --address ${ADDRESS} --port ${PORT} -e ${ENVIRONMENT} --pid-file ${PID} -d"

case "$1" in
  start)
    echo "Starting $APP_NAME passenger"
    echo "su - ${AS_USER} -c \"$CMD\""
    su - ${AS_USER} -c "$CMD"
    ;;
  stop)
    echo "Stopping $APP_NAME passenger"
    cd ${APP_ROOT}
    $PASSENGER stop -p ${PORT}
    ;;
  restart)
    echo "Restarting $APP_NAME passenger"
    cd ${APP_ROOT}
    $PASSENGER stop -p ${PORT}

    echo "su - ${AS_USER} -c \"$CMD\""
    su - ${AS_USER} -c "$CMD"
    ;;
  reload)
    echo "Reloading $APP_NAME passenger"
    cd ${APP_ROOT}
    $PASSENGER_CONFIG restart-app
    ;;
  *)
    echo "Usage: $0 start|stop|restart|reload" >&2
    exit 3
    ;;
esac
